# docs/watchdog.service.template
[Unit]
Description=WatchDog Daemon (daily Pulse at 21:00)
After=network-online.target

[Service]
Type=simple
User=watchdog
WorkingDirectory=/opt/watchdog
ExecStart=/opt/watchdog/.venv/bin/python -m watchdog.daemon
Restart=on-failure
RestartSec=5
EnvironmentFile=/opt/watchdog/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target